# root/etc/services.d/openssh-server/run
#!/usr/bin/with-contenv bash

USER_NAME=${USER_NAME:-ivonet.nl}

exec 2>&1 \
        s6-setuidgid root /usr/sbin/sshd -D -e -p 22
